# deploy/backup/nnp-backup.service
[Unit]
Description=NonNewtonian Physicists encrypted off-site backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/nnp-backup.env
# Adjust this path to wherever you put the repo/scripts on the server.
ExecStart=/opt/nnp/deploy/backup/nnp-backup.sh
# Be a good neighbour to the web app on a small box.
Nice=10
IOSchedulingClass=idle
